# alpine-consul/root/etc/services.d/consul/run
#!/usr/bin/with-contenv sh

if [ -z "$MASTER" ]
then
	next="-bootstrap"
else
	next="-join=$MASTER"
fi

exec consul agent -config-dir /etc/consul.d/bootstrap $next -dc $DATACENTER;